# target/linux/mediatek/mt7981/base-files/etc/init.d/telnetd
#!/bin/sh /etc/rc.common
#
START=50

start() {
	/usr/sbin/telnetd -l /bin/login <> /dev/null 2>&1
}
stop() {
	killall telnetd
}